# contrib/systemd/apm2-worker.service
[Unit]
Description=APM2 FAC Worker
Documentation=https://github.com/guardian-intelligence/apm2
After=network-online.target apm2-daemon.service
Wants=apm2-daemon.service
Requires=apm2-daemon.service

[Service]
# TCK-00600: Type=notify — worker sends READY=1 after broker connection,
# WATCHDOG=1 periodically. systemd restarts if watchdog is not pinged.
Type=notify
# Resource limits derived from LaneProfileV1 defaults
# CPUQuota=200%
# MemoryMax=25769803776
# TasksMax=1536
# IOWeight=100
# TimeoutStartSec=600
# RuntimeMaxSec=1800
# KillMode=control-group
CPUQuota=200%
MemoryMax=25769803776
TasksMax=1536
IOWeight=100
TimeoutStartSec=600
RuntimeMaxSec=1800
KillMode=control-group
# Adjust this path to match your deployment layout if different.
ExecStart=%h/.cargo/bin/apm2 fac worker --poll-interval 10
Restart=always
RestartSec=5
WatchdogSec=300
DynamicUser=yes
StateDirectory=apm2
CacheDirectory=apm2
LogsDirectory=apm2
CapabilityBoundingSet=
AmbientCapabilities=
NoNewPrivileges=yes
WorkingDirectory=/var/lib/apm2
Environment=APM2_HOME=/var/lib/apm2
Environment=RUST_LOG=info
LoadCredential=gh-token:/var/lib/apm2/private/creds/gh-token
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/apm2
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
