# contrib/systemd/apm2-daemon.service
# TCK-00595: System-mode service unit (multi-user.target).
# For user-mode, see contrib/systemd/user/apm2-daemon.service.
[Unit]
Description=APM2 Daemon — Forge Admission Cycle
Documentation=https://github.com/guardian-intelligence/apm2
After=network-online.target
Wants=network-online.target
Requires=apm2-daemon.socket

[Service]
# TCK-00600: Type=notify — daemon sends READY=1 after socket bind, WATCHDOG=1
# periodically. systemd restarts if watchdog is not pinged within WatchdogSec.
Type=notify
ExecStart=/usr/local/bin/apm2 daemon --no-daemon
ExecStop=/usr/local/bin/apm2 kill
Restart=always
RestartSec=5
WatchdogSec=300
DynamicUser=yes
StateDirectory=apm2
CacheDirectory=apm2
LogsDirectory=apm2
CapabilityBoundingSet=
AmbientCapabilities=
NoNewPrivileges=yes
WorkingDirectory=/var/lib/apm2
Environment=APM2_HOME=/var/lib/apm2
Environment=RUST_LOG=info
LoadCredential=gh-token:/var/lib/apm2/private/creds/gh-token
Sockets=apm2-daemon.socket
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/apm2
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
