# contrib/systemd/apm2-worker.service
[Unit]
Description=APM2 FAC Worker
Documentation=https://github.com/guardian-intelligence/apm2
After=network-online.target apm2-daemon.service
Wants=apm2-daemon.service
Requires=apm2-daemon.service

[Service]
Type=simple
# Adjust this path to match your deployment layout if different.
ExecStart=%h/.cargo/bin/apm2 fac worker --poll-interval 10
Restart=always
RestartSec=5
WatchdogSec=300
DynamicUser=yes
StateDirectory=apm2
CacheDirectory=apm2
LogsDirectory=apm2
CapabilityBoundingSet=
AmbientCapabilities=
NoNewPrivileges=yes
WorkingDirectory=/var/lib/apm2
Environment=APM2_HOME=/var/lib/apm2
Environment=RUST_LOG=info
LoadCredential=gh-token:/var/lib/apm2/private/creds/gh-token
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/apm2
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
