# contrib/systemd/user/apm2-daemon.service
[Unit]
Description=APM2 Daemon — Forge Admission Cycle broker
Documentation=https://github.com/guardian-intelligence/apm2
After=network-online.target
Wants=network-online.target
Requires=apm2-daemon.socket

[Service]
Type=simple
ExecStart=/usr/local/bin/apm2 daemon --no-daemon
ExecStop=/usr/local/bin/apm2 kill
Restart=always
RestartSec=5
WatchdogSec=300
WorkingDirectory=%h
Environment=APM2_HOME=%h/.apm2
Environment=RUST_LOG=info
Environment=XDG_RUNTIME_DIR=%t
LoadCredential=gh-token:%h/.apm2/private/creds/gh-token
Sockets=apm2-daemon.socket
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.apm2
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=default.target
