# contrib/systemd/user/apm2-daemon.service
# TCK-00595: systemd user service unit for apm2-daemon.
# Install via: apm2 daemon install --enable-linger
# The install command dynamically resolves ExecStart from the current binary.
# This reference file uses %h/.cargo/bin/apm2 as the default cargo install path.
# TCK-00600: Type=notify — daemon sends READY=1 after socket bind, WATCHDOG=1
# periodically. systemd restarts if watchdog is not pinged within WatchdogSec.
[Unit]
Description=APM2 Daemon — Forge Admission Cycle
After=network-online.target
Wants=network-online.target
Requires=apm2-daemon.socket

[Service]
Type=notify
ExecStart=%h/.cargo/bin/apm2 daemon --no-daemon
ExecStop=%h/.cargo/bin/apm2 kill
Restart=always
RestartSec=5
WatchdogSec=300
WorkingDirectory=%h/.apm2
Environment=APM2_HOME=%h/.apm2
Environment=RUST_LOG=info
Environment=XDG_RUNTIME_DIR=%t
# GITHUB_TOKEN is inherited from user session or read via LoadCredential.
# Never persist tokens in unit files (security policy).
LoadCredential=gh-token:%h/.apm2/private/creds/gh-token
Sockets=apm2-daemon.socket
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.apm2
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=default.target
